# systemd/palworld-backup-weekly.timer
[Unit]
Description=Weekly Palworld backup
[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true
[Install]
WantedBy=timers.target
