# systemd/palworld-public-status.timer
[Unit]
Description=Refresh static Palworld status
[Timer]
OnBootSec=30s
OnUnitActiveSec=15s
AccuracySec=2s
Persistent=true
[Install]
WantedBy=timers.target
